# package/wlan_10_2/qca-spectral/files/drm.init
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/usr/sbin/drm

start_service()
{
        procd_open_instance
        procd_set_param command "$PROG"
        procd_set_param respawn
        procd_close_instance
}

stop_service()
{
	start-stop-daemon -K -x "$PROG" -s SIGKILL > /dev/null
}

